# fieldrec/deploy/fieldrec-web.service
[Unit]
Description=Field recording web controller
After=network-online.target
Wants=audio-sync.service

[Service]
User=rec
WorkingDirectory=/opt/fieldrec
ExecStart=/opt/fieldrec/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fieldrec

[Install]
WantedBy=multi-user.target
